# deploy/launchd/com.egp.pg-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Track C — always-on SSH tunnel to PRODUCTION Postgres.
  __REPO_ROOT__ is substituted by scripts/install_launchd.sh at install time.
  KeepAlive restarts the tunnel if SSH drops (paired with ExitOnForwardFailure).
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.egp.pg-tunnel</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>exec "__REPO_ROOT__/scripts/run_remote_crawl.sh" tunnel</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/egp/tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/egp/tunnel.log</string>
  </dict>
</plist>
